# server/packaging/pi-manager.service
[Unit]
Description=pi-manager daemon
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi-manager
Group=pi-manager
ExecStart=/usr/local/bin/pi-manager --state /var/lib/pi-manager/state.json --addr 127.0.0.1:8080
Restart=on-failure
RestartSec=5
ProtectSystem=full
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
